# ENV-Chander-India-Tunnel/etc/update-motd.d/20-tunnel-logo
#!/bin/sh

# рҹ…І рҹ…· рҹ…° рҹ…Ҫ рҹ…і рҹ…ҙ рҹҶҒ
# -------------------
# Message of the day
#
# Location: /etc/update-motd.d folder
# Note: The file must starts with a number-name pattern, and
#       should NOT have an extension such as .sh or .bash etc.
# Example: 10-motd

# IMPORTANT: After dropping this file at the destination, run
#            sudo service ssh restart

{
    cat << "EOF"

в–—в–„в–„в–„в––в–—в––  в–—в––в–—в–„в–„в–„ в–—в–„в–„в–„в–– в–—в–„в––    в–—в–„в–„в–„в––в–—в–– в–—в––в–—в––  в–—в––в–—в––  в–—в––в–—в–„в–„в–„в––в–—в––
  в–Ҳ  в–җв–ӣв–ҡв––в–җв–Ңв–җв–Ң  в–Ҳ  в–Ҳ  в–җв–Ң в–җв–Ң     в–Ҳ  в–җв–Ң в–җв–Ңв–җв–ӣв–ҡв––в–җв–Ңв–җв–ӣв–ҡв––в–җв–Ңв–җв–Ң   в–җв–Ң
  в–Ҳ  в–җв–Ң в–қв–ңв–Ңв–җв–Ң  в–Ҳ  в–Ҳ  в–җв–ӣв–Җв–ңв–Ң     в–Ҳ  в–җв–Ң в–җв–Ңв–җв–Ң в–қв–ңв–Ңв–җв–Ң в–қв–ңв–Ңв–җв–ӣв–Җв–Җв–ҳв–җв–Ң
в–—в–„в–Ҳв–„в––в–җв–Ң  в–җв–Ңв–җв–ҷв–„в–„в–Җв–—в–„в–Ҳв–„в––в–җв–Ң в–җв–Ң     в–Ҳ  в–қв–ҡв–„в–һв–ҳв–җв–Ң  в–җв–Ңв–җв–Ң  в–җв–Ңв–җв–ҷв–„в–„в––в–җв–ҷв–„в–„в––

EOF

    uptime --p | sed 's/^./ \u&/'

} | LANG=en_US.UTF-8 /usr/games/lolcat -f
